# systemd/vps-sec-digest.timer
[Unit]
Description=vps-sec — agenda o resumo diário

[Timer]
OnCalendar=*-*-* 08:00:00
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
