# systemd/vps-sec-audit.timer
[Unit]
Description=vps-sec — agenda a auditoria diária

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
